# deployment/deception-detector-frontend.service
[Unit]
Description=Deception Detector Frontend (Vue.js/Node.js)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/deception-detector/frontend
Environment="NODE_ENV=production"
Environment="PORT=8080"
ExecStart=/usr/bin/node /opt/deception-detector/frontend/server.js

# Restart settings
Restart=always
RestartSec=10
StartLimitInterval=400
StartLimitBurst=3

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

# Resource limits
MemoryLimit=2G
CPUQuota=200%

[Install]
WantedBy=multi-user.target
