# deployment/deception-detector-backend.service
[Unit]
Description=Deception Detector Backend (Flask API)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/deception-detector/backend
Environment="PATH=/opt/deception-detector/backend/venv/bin"
EnvironmentFile=/opt/deception-detector/backend/.env
ExecStart=/opt/deception-detector/backend/venv/bin/gunicorn \
    --workers 4 \
    --threads 2 \
    --timeout 300 \
    --bind 127.0.0.1:5000 \
    --access-logfile /var/log/deception-detector/backend-access.log \
    --error-logfile /var/log/deception-detector/backend-error.log \
    --log-level info \
    app:app

# Restart settings
Restart=always
RestartSec=10
StartLimitInterval=400
StartLimitBurst=3

# Security settings (relaxed for compatibility)
NoNewPrivileges=true

# Resource limits (adjust based on your server)
MemoryLimit=8G
CPUQuota=400%

[Install]
WantedBy=multi-user.target
